# Grape_Gen1_PSWS/PSWSsetup/xfer2repo.crontab
42 00 * * * /home/pi/PSWS/Scode/xfer2repo.sh >/home/pi/PSWS/Sstat/xferstat 2>&1
